# crontab
0 8 * * * /app/pykrotik.py --devices /devices.yaml --output /output

